# systemd_services/tk-hello-onetime.service
[Unit]
Description=Hello Service (Onetime)

[Service]
ExecStart=/usr/local/bin/tk-service-hello.sh

[Install]
WantedBy=multi-user.target
